# opensciencegrid/osgstorage-redirector/supervisord.d/cert-watch.conf
[program:auto-reload]
command=/usr/local/sbin/cert-watch.sh
autorestart=true
