# opensciencegrid/frontier-squid/supervisord.d/fronted-squid.conf
[program:frontier-squid]
command=/usr/local/sbin/entrypoint.sh
autorestart=true
