# opensciencegrid/osdf-collector/etc/supervisord.d/htcondor.conf
[program:htcondor]
command=/usr/sbin/condor_master -f
autorestart=true
startsecs=20
