# opensciencegrid/ospool-ganglia/05-gmetad.conf
[program:gmetad]
command=/usr/sbin/gmetad -d 1
autorestart=True
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

