# opensciencegrid/ospool-ganglia/05-php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm --allow-to-run-as-root --nodaemonize
autorestart=true
startsecs=5
redirect_stderr=true

