# opensciencegrid/koji-builder/kojid-supervisord.conf
[program:kojid]
command=/usr/sbin/kojid
