# opensciencegrid/koji-builder/kojid-supervisord.conf
[program:kojid]
command=/usr/sbin/kojid --fg --force-lock --verbose

[program:kojid-tail]
command=tail -f /var/log/kojid.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
