# opensciencegrid/central-collector/etc/supervisord.d/auto-reload.conf
[program:auto-reload]
command=/usr/local/sbin/auto-reload.sh
autorestart=true
