# opensciencegrid/central-collector/etc/supervisord.d/htcondor-ce-collector.conf
[program:htcondor-ce-collector]
command=/usr/share/condor-ce/condor_ce_startup collector -f
autorestart=true
startsecs=20
