# opensciencegrid/frontier-squid/supervisor-frontier-squid.conf
[program:frontier-squid]
command=/usr/sbin/squid --foreground
user=squid
autorestart=true
