# opensciencegrid/frontier-squid/supervisor-frontier-squid.conf
[program:frontier-squid]
environment=SQUID_START_ARGS="--foreground"
command=/usr/sbin/fn-local-squid.sh start
stopwaitsecs=45
user=squid
autorestart=true
