# opensciencegrid/frontier-squid/supervisor-frontier-squid.conf
[program:frontier-squid]
command=/usr/sbin/start-frontier-squid.sh
stopwaitsecs=90
user=squid
