# opensciencegrid/sssd/etc/supervisord.d/sssd.conf
[program:sssd]
command=/usr/sbin/sssd --interactive --logger=files
autorestart=true
