# opensciencegrid/gwms-factory/etc/supervisord.d/10_gwms_factory.conf

[program:http]
command=/usr/sbin/httpd -DFOREGROUND

[program:factory]
command=python /usr/sbin/glideinFactory.py /var/lib/gwms-factory/work-dir
user=gfactory
autorestart=true

