# opensciencegrid/ospool-display/supervisor.d/10-fifemon.conf
[program:fifemon]
command=/usr/bin/fifemon /etc/fifemon.cfg
user=root
autorestart=true
startsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

