# opensciencegrid/vo-frontend/frontend_scitoken_check_ces.cron

48 * * * *  root  /usr/bin/timeout 3590 /usr/sbin/frontend_scitoken_check_ces >> /var/log/frontend_scitoken_check_ces.log 2>&1

