# opensciencegrid/ospool-cm/10-htcondor.conf
[program:condor_master]
command=/usr/sbin/condor_master_wrapper
autorestart=True
startsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

