# opensciencegrid/ospool-cm/30-prometheus.conf
[program:prometheus]
command=/prometheus.py
autorestart=True
startsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

