# opensciencegrid/ospool-cm/ospool.cron

# prios
 */5 * * * *  root  /opt/ospool/update-prios >/dev/null 2>&1

# web
 */5 * * * *  root  cd /shared/web/overview && ./utils/cron-5m.sh >/dev/null 2>&1
  0  * * * *  root  cd /shared/web/overview && ./utils/cron-hourly.sh >/dev/null 2>&1

